# templates/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
user=mysql
command=/usr/sbin/mysqld
